# files/etc/services.d/10_dnsmasq/run
#!/command/execlineb -P
/usr/sbin/dnsmasq -d --conf-dir /etc/dnsmasq.d
