# files/etc/services.d/02_docker-gen/run
#!/command/execlineb -P
with-contenv
/usr/local/bin/docker-gen -watch -notify-output -notify "run-parts /etc/docker-gen/notify.d" /etc/hosts.dnsmasq.tmpl /etc/hosts.dnsmasq
